# supervisor/conf.d/supervisord-logs.conf
[program:nginx-access-log]
command=tail -f /var/log/nginx/access.log
redirect_stderr=true

[program:nginx-error-log]
command=tail -f /var/log/nginx/error.log
redirect_stderr=true
