# supervisor/conf.d/supervisord-nginx.conf
[program:nginx]
command=/etc/supervisor/commands/start-nginx.sh
priority=999
autostart=true
autorestart=true
redirect_stderr=true
